# systemd/.config/systemd/user/swww_change_wallpaper.service
[Unit]
Description=Change wallpaper using SWWW
After=swww.service

[Service]
ExecStart=/bin/bash /home/ngumbk/.config/swww/swww.sh

[Install]
WantedBy=default.target
